# roles/local/20-docker/templates/etc/systemd/system/docker.service.d/options.conf.j2
# {{ ansible_managed }}

[Service]
ExecStart=
ExecStart=/usr/bin/dockerd {{ docker__daemon_location_flags | join(" ") }} {{ docker__daemon_flags | join(" ") }}
