# roles/local/05-environment/templates/wol.service.j2
[Unit]
Description=Configure Wake On LAN

[Service]
Type=oneshot
ExecStart=/sbin/ethtool -s {{ ansible_default_ipv4.interface }} wol g

[Install]
WantedBy=basic.target
